# soc-exporter.service
[Unit]
Description=SOC Exporter — Wazuh event forwarder
Documentation=https://soc-api.nox5.com.br
After=network-online.target wazuh-manager.service
Wants=network-online.target

[Service]
Type=simple
User=soc-exporter
Group=soc-exporter
ExecStart=/usr/local/bin/soc-exporter start
Restart=always
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=soc-exporter

# Ensure config and data dirs are accessible
ReadWritePaths=/etc/soc-exporter /var/lib/soc-exporter /var/log/soc-exporter
ReadOnlyPaths=/var/ossec/logs/alerts

# Harden: drop unnecessary privileges
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
